# adsl_check/adsl_check
#!/bin/sh /etc/rc.common
START=50

# this script should be placed inside /etc/init.d/

start() {

adsl_check.sh &

}

stop() {

killall adsl_check.sh

}

